# files/etc/crontabs/root
*/4 * * * * bash /etc/xyw/login.sh
0 4 * * * /etc/init.d/youhua start
