# PATCH/files/etc/uci-defaults/98-pre_install
#!/bin/sh

[ -n "$(ls /etc/pre_install)" ] && opkg install /etc/pre_install/*.ipk --force-depends
rm -rf /etc/pre_install

uci set fstab.@global[0].check_fs=1
uci commit

exit 0