# files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
bash /etc/xyw/login.sh
iptables -F
iptables -X
iptables -Z

exit 0
